# systemd/aql-agent.service
[Unit]
Description=AQL Agent — AI assistant on Raspberry Pi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ibpi
WorkingDirectory=/home/ibpi/aql-agent
EnvironmentFile=/home/ibpi/aql-agent/.env
ExecStart=/usr/bin/node /home/ibpi/aql-agent/dist/index.js
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aql-agent

[Install]
WantedBy=multi-user.target
